# vagrant/postgres.service
[Unit]
Description=Start PG

[Service]
Type=forking
User=vagrant
ExecStart=/usr/local/pgsql/bin/pg_ctl -D /media/data/pg_data/ start

# Due to PostgreSQL's use of shared memory, OOM killer is often overzealous in
# killing Postgres, so adjust it downward
OOMScoreAdjust=-200

# Additional security-related features
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
